# aconfmgr/.config/aconfmgr/files/etc/xdg/autostart/xcape.desktop
[Desktop Entry]
Name=xcape
Description=Use CTRL as ESC when pressed alone
Exec=xcape
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
