# aconfmgr/.config/aconfmgr/files/etc/systemd/system/undervolt.service
[Unit]
Description=Intel Undervolt Service
After=multi-user.target suspend.target hibernate.target hybrid-sleep.target

[Service]
Type=oneshot
ExecStart=/usr/bin/undervolt -t 97 --core -107 --cache -107 --gpu -83 --uncore -107

[Install]
WantedBy=multi-user.target suspend.target hibernate.target hybrid-sleep.target
