# aconfmgr/.config/aconfmgr/files/etc/profile.d/libva-driver.sh
export LIBVA_DRIVER_NAME="iHD"

# Fix MPV.
export MESA_LOADER_DRIVER_OVERRIDE="i965"
