# aconfmgr/.config/aconfmgr/files/etc/profile.d/wayland.sh
export CLUTTER_BACKEND=wayland
export SDL_VIDEODRIVER=wayland
export QT_QPA_PLATFORM=wayland-egl

export MOZ_ENABLE_WAYLAND=1
