# aconfmgr/.config/aconfmgr/files/etc/profile.d/wayland-toolkits.sh
export CLUTTER_BACKEND=wayland
export SDL_VIDEODRIVER=wayland
export QT_QPA_PLATFORM=wayland-egl
