# gnome/.config/autostart/jbl-setup.desktop
[Desktop Entry]
Type=Application
Name=JBL Setup
Exec=jbl-setup.sh
Terminal=false
