# gnome/.config/autostart/startup.desktop
[Desktop Entry]
Type=Application
Name=Startup
Exec=startup.sh
