# etc/systemd/system/undervolt.service
[Unit]
Description=Undervolt the machine

[Service]
Type=oneshot
ExecStart=/home/andrea/.bin/undervolt

[Install]
WantedBy=multi-user.target
